# linux/etc/systemd/system/btrfs-scrub@.service
[Unit]
Description=Scrub btrfs filesystems

[Service]
Type=oneshot
ExecStart=/home/vacation/bin/btrfs-scrub.sh %i
NotifyAccess=main
