# dotfiles/.config/systemd/user/docker-prune.service
[Unit]
Description=Prune docker system

[Service]
Type=simple
ExecStart=echo y | docker system prune
Restart=no
RestartSec=10s
