# dotfiles/.config/systemd/user/srit.service
[Unit]
Description=Apply keystroke speed tweak

[Service]
Type=simple
ExecStart=/usr/bin/xset r rate 140 70

[Install]
WantedBy=graphical.target
