# dotfiles/.config/systemd/user/xset-rate.service
[Unit]
Description=Set keystroke speed

[Service]
Type=simple
ExecStart=/usr/bin/xset r rate 140 70

[Install]
WantedBy=graphical.target
