# dotfiles/.config/systemd/user/xset-rate.timer
[Unit]
Description=Set keystroke speed timer

[Timer]
OnUnitActiveSec=15
AccuracySec=1us
Persistent=false

[Install]
WantedBy=timers.target
