# linux/home/.config/systemd/user/guake.service
[Unit]
Description=guake dropdown terminal process

[Service]
Type=simple
ExecStart=/usr/bin/guake

[Install]
WantedBy=default.target
