# linux/home/.config/systemd/user/tilda.service
[Unit]
Description=Tilda dropdown terminal process

[Service]
Type=simple
ExecStart=/usr/bin/tilda

[Install]
WantedBy=default.target
