# linux/home/.config/systemd/user/yakuake.service
[Unit]
Description=yakuake dropdown terminal process

[Service]
Type=simple
ExecStart=/usr/bin/yakuake

[Install]
WantedBy=default.target
