# etc/systemd/system/btrfs-snapshot.service
[Unit]
Description=btrfs snapshot system and delete old snapshots

[Service]
Type=oneshot
ExecStart=/home/vacation/bin/system-snapshot
User=root
