# etc/systemd/system/system-snapshot.service
[Unit]
Description=btrfs snapshot system and delete old snapshots

[Service]
Type=oneshot
ExecStart=%h/bin/system-snapshot
User=root
